# systemd/openvswitch-persistence.service.in
[Unit]
Description=Apply Open vSwitch persistence configuration.
Before=openvswitch.service
After=openvswitch-nonetwork.service
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
ExecStart=@DESTDIR@@DATA_PREFIX@/openvswitch-persistence/apply.py
ExecStop=@DESTDIR@@DATA_PREFIX@/openvswitch-persistence/apply.py

[Install]
WantedBy=openvswitch.service
